# automation/supervisord/supervisord.conf
[supervisord]
nodaemon=false

[program:server]
command=/opt/game/server
environment=DISCORDTOKEN="%(ENV_DISCORDTOKEN)s"
autostart=true
autorestart=true
redirect_stderr=true