# automation/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=./server.exe
#environment=DISCORD_TOKEN="your_discord_token"
autostart=true
autorestart=true
redirect_stderr=true