# automation/supervisord/supervisord.conf
[supervisord]
nodaemon=false

[program:nexus]
directory=/opt/game/
command=./server
environment=DISCORDTOKEN="%(ENV_DISCORDTOKEN)s"
autostart=true
autorestart=true
redirect_stderr=true